# rc.local
#!/bin/bash

chown shadow /sys/class/leds/asus::kbd_backlight/brightness
